# curator4-600467809.conf
[program:curator4-600467809]
command=python3 /home/curator4/curator4-600467809/main.py
autostart=true
autorestart=true
stderr_logfile=/home/curator4/curator4-600467809/test.err.log
stdout_logfile=/home/curator4/curator4-600467809/test.out.log